# static/terminals.slice
[Unit]
Description=Terminal cgroups
Delegate=yes

[Install]
WantedBy=default.target
