# static/terminals.slice
[Unit]
Description=Terminal cgroups

[Install]
WantedBy=default.target
